# deploy/systemd/legs-dog.service
# deploy/systemd/legs-dog.service
[Unit]
Description=legs dog runtime (sim/real + low-latency loop)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=legs
WorkingDirectory=/opt/legs
EnvironmentFile=/opt/legs/deploy/env/legs.env
ExecStart=/usr/bin/docker compose -f /opt/legs/deploy/docker/docker-compose.planA.yml up --remove-orphans
ExecStop=/usr/bin/docker compose -f /opt/legs/deploy/docker/docker-compose.planA.yml down
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
